# crontab
0 */3 * * * bash /home/scripts/nod32-mirror/nod32-mirror.sh --update
