# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
autorestart=true

[program:cron]
directory=/etc/cron.d/
command=/usr/sbin/cron -f
user=root
autorestart=true
